# c/ppp.service.in
[Unit]
Description=PPP Server

[Service]
ExecStart=@PREFIX@/libexec/ppp_start.sh
Restart=no
Type=exec

[Install]
WantedBy=multi-user.target
